# infra/portscan-client.service
[Unit]
Description=Port Scanner Client Agent
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/portscan-client
ExecStart=/opt/portscan-client/venv/bin/python client_agent.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/portscan-client

[Install]
WantedBy=multi-user.target
